# launchd/com.freedom33.truthsocial.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.freedom33.truthsocial.monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/freedom33/Documents/New_project/scripts/truthsocial_watcher.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/freedom33/Documents/New_project</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/freedom33/Documents/New_project/logs/truthsocial_launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/freedom33/Documents/New_project/logs/truthsocial_launchd.err.log</string>
</dict>
</plist>
